# systemd/nearby.service
[Unit]
Description=Nearby Service

[Service]
ExecStart=/usr/bin/env nearby

[Install]
WantedBy=default.target
